# flog.cron
*/15 * * * * /home/pi/flog.sh >> flog.log
